# auto-restart/85-remove-solokey.rules
ACTION=="remove", SUBSYSTEM=="usb", ENV{PRODUCT}=="SOLOKEY_PRODUCT_ID", RUN+="/usr/local/bin/solo"
