# auto-restart/85-remove-usb.rules
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/usr/local/bin/solo"
